# etc/systemd/system/getty@.service.d/override.conf
[Service]
Type=simple
ExecStart=
ExecStart=-/sbin/agetty --autologin hch --noclear %I $TERM
